# aw-server/dlp-integrations/webhook-sender.service
[Unit]
Description=AWatch DLP Webhook Sender
After=network-online.target

[Service]
Type=oneshot
User=activitywatch
Group=activitywatch
WorkingDirectory=/opt/activitywatch/dlp-integrations
ExecStart=/opt/activitywatch/dlp-integrations/.venv/bin/python /opt/activitywatch/dlp-integrations/webhook_sender.py

[Install]
WantedBy=multi-user.target
