# aw-server/aw-worktime-ui-bridge.service
[Unit]
Description=AW Worktime UI bridge (sessions -> afk/window)
After=network-online.target activitywatch-server.service
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=120

[Service]
Type=simple
Environment=AW_SERVER_URL=http://127.0.0.1:5600
Environment=AW_WORKTIME_HOST=SHARKON2025
ExecStart=/usr/bin/python3 /usr/local/bin/aw-worktime-ui-bridge.py
Restart=on-failure
RestartSec=10
User=activitywatch
Group=activitywatch
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aw-worktime-ui-bridge

[Install]
WantedBy=multi-user.target
