# aw-server/dlp-compliance/report-scheduler.service
[Unit]
Description=AWatch DLP 152-FZ compliance report generator
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/activitywatch/aw-server.env
WorkingDirectory=/opt/activitywatch/dlp-compliance
ExecStart=/opt/activitywatch/dlp-compliance/.venv/bin/python /opt/activitywatch/dlp-compliance/compliance_scheduler.py
User=activitywatch
Group=activitywatch
