# aw-server/dlp-policy-engine/dlp-policy-engine.service
[Unit]
Description=AW DLP Policy Engine
After=network.target activitywatch-server.service
Wants=activitywatch-server.service

[Service]
Type=simple
EnvironmentFile=/etc/activitywatch/aw-server.env
WorkingDirectory=/opt/activitywatch/dlp-policy-engine
ExecStart=/opt/activitywatch/dlp-policy-engine/.venv/bin/uvicorn policy_service:app --host ${AW_DLP_POLICY_ENGINE_BIND_HOST} --port ${AW_DLP_POLICY_ENGINE_PORT}
Restart=on-failure
RestartSec=5
StartLimitBurst=3
StartLimitIntervalSec=60
User=activitywatch
Group=activitywatch
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aw-dlp-policy-engine

[Install]
WantedBy=multi-user.target
