# pfsense/pfsense-aw-poller.service
[Unit]
Description=pfSense ActivityWatch poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/aw-pfsense/pfsense-aw-poller.py --config /etc/aw-pfsense/poller.json
Restart=always
RestartSec=10
User=root
Group=root

[Install]
WantedBy=multi-user.target
