# install-kit-awindows-20260427-211240/aw-server/aw-rus-healthd.timer
[Unit]
Description=Run AW-RUS unified health orchestrator every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=aw-rus-healthd.service
Persistent=true

[Install]
WantedBy=timers.target
